# build/supervisor/supervisord.conf
[supervisor]
loglevel=info
redirect_stderr=true

[program:laravel_schedule_work]
directory=/var/www/html
command=php artisan schedule:work
priority=999
numprocs=1
startsecs=0
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/laravel_schedule.log

[program:laravel_queue_normal]
directory=/var/www/html
command=php artisan queue:work --queue=normal
process_name=%(program_name)s_%(process_num)s
priority=998
numprocs=2
startsecs=0
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/laravel_queue_normal.log